# dot_config/systemd/user/app-foot@.service.d/override.conf
[Service]
# zellij dumps session layout/metadata on SIGTERM; default 90s can be too short
# for large sessions, causing SIGKILL mid-write and unresurrectable state.
TimeoutStopSec=300
