# dot_config/systemd/user/foot-server.service.d/override.conf
[Service]
# zellij runs under foot-server's cgroup. On shutdown, give zellij time to
# serialize session layout/metadata before SIGKILL. Default 90s was too short
# and caused mid-write kills (see main session metadata loss, 2026-04-18).
TimeoutStopSec=300
